# profile.d/aliases.sh
alias ls='ls -hF --color=tty'
alias ll='ls -Al'
alias l='ls -l'
alias grep='grep --color=tty'
