# profile.d/history.sh
shopt -s histappend

export HISTCONTROL=ignoreboth
export HISTFILESIZE=1000
export HISTIGNORE='ls:l:ll'
export HISTTIMEFORMAT='[%F %T] '
